# deployment/systemd/culifeed-processor.timer
[Unit]
Description=CuliFeed Daily Content Processing Timer
Documentation=https://github.com/chiplonton/culifeed
Requires=culifeed-processor.service
PartOf=culifeed.target

[Timer]
# Run daily at 8:00 AM local time
OnCalendar=*-*-* 08:00:00
# Run 2 minutes after system boot (for initial setup)
OnBootSec=2min
# If the system was off when the timer should have run, run it immediately when it starts
Persistent=true
# Add some randomization (±15 minutes) to avoid thundering herd
RandomizedDelaySec=900

[Install]
WantedBy=timers.target