# deployment/systemd/culifeed-bot.service
[Unit]
Description=CuliFeed Telegram Bot Service
Documentation=https://github.com/chiplonton/culifeed
After=network-online.target
Wants=network-online.target
PartOf=culifeed.target

[Service]
Type=exec
User=culifeed
Group=culifeed
WorkingDirectory=/opt/culifeed
Environment="CULIFEED_CONFIG=/etc/culifeed/config.yaml"
EnvironmentFile=/etc/culifeed/.env
ExecStart=/opt/culifeed/venv/bin/python /opt/culifeed/run_bot.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=512M
TasksMax=200

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/culifeed /var/lib/culifeed /tmp
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=culifeed-bot

# Health monitoring
WatchdogSec=30

[Install]
WantedBy=multi-user.target